# qeli/debian/qeli.service
[Unit]
Description=Qeli VPN Server
Documentation=https://qeli.dev
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
# Accept sd_notify STATUS= from the main (supervisor) process so `systemctl status`
# shows the running qeli version + profile count. Type stays `simple` (no READY
# handshake); a missing socket makes the notify a harmless no-op.
NotifyAccess=main
User=qeli
Group=qeli
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=full
# ProtectSystem=full mounts /etc read-only, but the server writes inside
# /etc/qeli at runtime: it generates the per-profile identity key
# (/etc/qeli/identity/<name>.key) on first start, and the web panel / control
# commands persist user changes to the users file (default /etc/qeli/users.conf).
# Without this exception those writes fail with EROFS — the profile fails to bind
# on a fresh install, and `disable-user` / panel edits silently don't persist.
ReadWritePaths=/etc/qeli
ProtectHome=true
PrivateTmp=true
StateDirectory=qeli
RuntimeDirectory=qeli
LogsDirectory=qeli

ExecStart=/usr/bin/qeli server --config /etc/qeli/server.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
LimitNPROC=1024

[Install]
WantedBy=multi-user.target